# tags/CR/pbs/submit_host/etc/cron.daily/svn-pbs-scripts
#!/bin/sh
cd /home/pbsuser/
#svn checkout --force http://hpc-cc.googlecode.com/svn/tags/CR/pbs/user/ .
svn --force update
cd -
exit 0
